# root.netfs/etc/inittab
# inittab pour linux
id:2:initdefault:
1:12345:respawn:/bin/mingetty tty1 --autologin axel
